# systemd/workproof-admin.service
[Unit]
Description=WorkProof Chain Admin Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/workproof-chain-v2
ExecStart=/usr/bin/python3 /home/ubuntu/workproof-chain-v2/src/admin_bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
